# systemd/system/openrgb-off.service
[Unit]
Description=Force RGB off at boot with OpenRGB
After=systemd-udev-settle.service
Wants=systemd-udev-settle.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/openrgb_force_off.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
